# deployment/systemd/homeos-auto-update.timer
# HomeOS Auto-Update Timer Template
#
# NOTE: This is a TEMPLATE/EXAMPLE file only and is NOT used directly.
# The actual timer file is generated dynamically by setup-auto-update.sh
# and installed to /etc/systemd/system/homeos-auto-update.timer
#
# This file shows the structure of the generated timer file for reference.

[Unit]
Description=HomeOS Auto-Update Timer
Requires=homeos-auto-update.service

[Timer]
# Check for updates every 10 minutes
OnBootSec=2min
OnUnitActiveSec=10min

# Run even if last run was missed (e.g., server was off)
Persistent=true

[Install]
WantedBy=timers.target
